# build/silk/service/rwflowpack/run
/usr/sbin/rwflowpack --pack-interfaces --site-config-file /config/silk.conf --sensor-configuration /config/sensor.conf --root-directory=/data/ --log-directory=/logs/
